# releases/beta-4.6.3/systemd/examples/halo-fire-standby.service.example
[Unit]
Description=Halo FIRE Standby NAS-Safe Report

[Service]
Type=oneshot
ExecStart=/usr/local/bin/halo-fire
