# releases/beta-4.6.3/systemd/examples/halo-fire-standby.timer.example
[Unit]
Description=Run Halo FIRE Standby NAS-Safe every 30 minutes

[Timer]
OnBootSec=3min
OnUnitActiveSec=30min
Unit=halo-fire-standby.service
Persistent=true

[Install]
WantedBy=timers.target
