# releases/beta-4.6.3/systemd/examples/halo-log-rotate.service.example
[Unit]
Description=Halo HomeLab Log Rotation

[Service]
Type=oneshot
ExecStart=/usr/local/bin/halo-log-rotate
